# telme10.socket
[Unit]
Description=<<</>> telnet server socket

[Socket]
ListenStream=23

[Install]
WantedBy=sockets.target
